# prometheus/discover-targets.service
# Copyright (c) Advanced Micro Devices, Inc. All rights reserved.
#
# SPDX-License-Identifier: MIT

[Unit]
Description=Discover Prometheus targets via Avahi/mDNS
After=avahi-daemon.service
Wants=avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/discover-targets.sh --deploy
User=root
